# src/main/scripts/HelloWorld.service
[Unit]
Description=Service run helloworld.sh

[Install]
WantedBy = multi-user.target

[Service]
Type=forking
ExecStart = /bin/helloworld.sh
